# templates/default/midas.upstart.erb
description "Start the midas node"

start on startup
stop on shutdown

respawn
respawn limit 10 5

setuid <%= @app_user %>

chdir <%= @working_dir %>

exec \
node app.js --production >> <%= @working_dir %>/midas.sys.log 2>&1
