# templates/default/midas.upstart.erb
description "Start the midas node"

start on startup
stop on shutdown

respawn
respawn limit 10 5

env NODE_ENV=<%= @app_environment %>

env NEW_RELIC_APP_NAME=<%= @newrelic_appname %>
env NEW_RELIC_LICENSE_KEY=<%= @newrelic_licensekey %>

env MYUSA_CLIENT_ID=<%= @myusa_client_id %>
env MYUSA_CLIENT_SECRET=<%= @myusa_client_secret %>
env MYUSA_CALLBACK_URL="https://<%= @app_host %>/api/auth/callback/myusa"

env LINKEDIN_CLIENT_ID=<%= @linkedin_client_id %>
env LINKEDIN_CLIENT_SECRET=<%= @linkedin_client_secret %>
env LINKEDIN_CALLBACK_URL="https://<%= @app_host %>/api/auth/callback/linkedin"

setuid <%= @app_user %>

chdir <%= @working_dir %>

exec \
node app.js >> <%= @working_dir %>/midas.sys.log 2>&1
