# templates/default/midas.upstart.erb
description "Start the midas node"

start on startup
stop on shutdown

respawn
respawn limit 10 5

env NEW_RELIC_APP_NAME=<%= @newrelic_appname %>
env NEW_RELIC_LICENSE_KEY=<%= @newrelic_licensekey %>

setuid <%= @app_user %>

chdir <%= @working_dir %>

exec \
node app.js >> <%= @working_dir %>/midas.sys.log 2>&1
